# deploy/file-loading-api.service
[Unit]
Description=Selcomm File Loading API v4
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/www/api/v4/file_loading_api
ExecStart=/var/www/api/v4/file_loading_api/FileLoading
Restart=on-failure
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=file_loading_api

User=weblocal
Group=webusers

Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5140
Environment=SELCOMM_CONFIG_PATH=/etc/selcomm/appsettings.shared.json
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=LD_LIBRARY_PATH=/opt/ibm/db2/clidriver/lib

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true

ReadWritePaths=/var/log/file-loading-api
ReadWritePaths=/var/www/api/v4/file_loading_api/logs
ReadWritePaths=/var/www

TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
